# maxtaylordavi.es.service
[Unit]
Description=maxtaylordavi.es

[Service]
Type=simple
Restart=always
RestartSec=5s
ExecStart=/home/pi/code/maxtaylordavi.es/maxtaylordavi.es
WorkingDirectory=/home/pi/code/maxtaylordavi.es/

[Install]
WantedBy=multi-user.target